# deploy/lxc/luma-deploy.service
[Unit]
Description=Pull luma repo and re-up the docker compose stack
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
Environment=LUMA_DIR=/opt/luma
Environment=LUMA_BRANCH=main
Environment=APP_PORT=3000
WorkingDirectory=/opt/luma
ExecStart=/bin/bash /opt/luma/deploy/lxc/luma-deploy.sh

[Install]
WantedBy=multi-user.target
