# deploy/luma-zoho-po-sync.timer
[Unit]
Description=Daily Luma Zoho purchase-order sync at 03:59 US Eastern

[Timer]
# Fire at 03:59 on the LXC host clock. The LXC must use America/New_York
# (see deploy/README-zoho-po-sync.md — timedatectl set-timezone).
OnCalendar=*-*-* 03:59:00
Persistent=true
AccuracySec=1min
Unit=luma-zoho-po-sync.service

[Install]
WantedBy=timers.target
