# deploy/luma-zoho-po-sync.timer
[Unit]
Description=Daily Luma Zoho purchase-order sync at 03:59 US Eastern

[Timer]
# LXC host runs UTC, but ops are on US Eastern. Fire at 03:59 America/New_York
# (EST/EDT with DST) regardless of the system clock timezone.
Timezone=America/New_York
OnCalendar=*-*-* 03:59:00
Persistent=true
AccuracySec=1min
Unit=luma-zoho-po-sync.service

[Install]
WantedBy=timers.target
